# launchd/com.sakesaakstra.btc-data-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sakesaakstra.btc-data-update</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/sakesaakstra/Desktop/crypto_hartrace/scripts/btc_data_runner.sh</string>
    </array>
    
    <!-- Dagelijks om 22:00 UTC (5 min vóór trader-job) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>22</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/sakesaakstra/Desktop/crypto_hartrace/logs/launchd_data.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/sakesaakstra/Desktop/crypto_hartrace/logs/launchd_data.err.log</string>
    
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
